# deploy/pve/botnote-dns.service
[Unit]
Description=IMTégrale local LAN DNS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botnote-dns
Group=botnote-dns
ExecStart=/usr/sbin/dnsmasq --keep-in-foreground --conf-file=/etc/botnote-dnsmasq.conf
Restart=on-failure
RestartSec=2
TimeoutStopSec=10
MemoryMax=64M
TasksMax=32
UMask=0077

AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=192.168.50.0/24
IPAddressAllow=100.64.0.0/10
IPAddressAllow=1.1.1.1/32

[Install]
WantedBy=multi-user.target
