# deploy/botnote-backup.timer
[Unit]
Description=Daily IMTégrale database backup

[Timer]
OnCalendar=*-*-* 03:40:00
RandomizedDelaySec=10min
Persistent=true
Unit=botnote-backup.service

[Install]
WantedBy=timers.target
