# deploy/botnote-backup.service
[Unit]
Description=Backup the IMTégrale PostgreSQL database
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
User=botnote
Group=botnote
Environment=BOTNOTE_BACKUP_DIR=/var/backups/botnote
Environment=BOTNOTE_BACKUP_RETENTION_DAYS=30
Environment=BOTNOTE_BACKUP_AGE_RECIPIENT_FILE=/etc/botnote/backup-age-recipient
ExecStart=/opt/botnote/current/deploy/backup.sh
MemoryMax=512M
TasksMax=64
UMask=0077

CapabilityBoundingSet=
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
ReadWritePaths=/var/backups/botnote
ReadOnlyPaths=/etc/botnote/backup-age-recipient
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
